# deployment/systemd.service
[Unit]
Description=Ledger Quick Note

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/path/to/lazy-finance
ExecStart=/path/to/lazy-finance/lazy-finance -f /path/to/ledger/journal.txt -w /path/to/ledger

[Install]
WantedBy=multi-user.target
